# etc/systemd/system/watch-qwen-serve.timer
[Unit]
Description=Daily check for Qwen Code serve/daemon/ACP changes

[Timer]
OnCalendar=*-*-* 10:00:00
RandomizedDelaySec=1800
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
